# roles/nvidia_container_toolkit/files/70-nvidia-ctk.rules
ACTION=="add|bind", KERNEL=="nvidia", RUN+="/usr/bin/nvidia-ctk cdi generate --output=/var/run/cdi/nvidia.yaml"